# optimizer/optimizer-daemon.service
[Unit]
Description=Maestro Optimizer Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/optimizer-daemon
ExecReload=/bin/kill -USR2 $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
PIDFile=/run/optimizer.pid

# Dedicated unprivileged user/group
User=maestro
Group=maestro

# systemd will create /run/optimizer with this owner
RuntimeDirectory=optimizer
RuntimeDirectoryMode=0750

# Log file is created beforehand and owned by maestro
StandardOutput=append:/var/log/maestro.log
StandardError=append:/var/log/maestro.log

[Install]
WantedBy=multi-user.target
